# docs/supervisor/gdc.conf
[program:gdc]
command=/var/www/bin/event_loop
autostart=true
autorestart=true
stderr_logfile=/var/log/gdc.err.log
stdout_logfile=/var/log/gdc.out.log
